# roles/laptop/files/kb-layout.rules
ATTRS{idVendor}=="413c", ATTRS{idProduct}=="2003", RUN+="/home/gianluca/bin/kb-layout us"
